# templates/default/epipe.service.erb
[Unit]
Description = Oozie Server
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.epipe.user %>

PIDFile = /tmp/epipe.pid
ExecStart = <%= node.epipe.home %>/bin/epipe-start.sh
ExecStop = <%= node.epipe.home %>/bin/epipe-stop.sh
LimitNOFILE = 10000

Type = forking

[Install]
WantedBy = multi-user.target
