# templates/default/epipe.service.erb
[Unit]
Description = ePipe Server
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.epipe.user %>

PIDFile = <%= node.epipe.pid_file %>
ExecStart = <%= node.epipe.home %>/bin/start-epipe.sh
ExecStop = <%= node.epipe.home %>/bin/stop-epipe.sh
LimitNOFILE = 20000

Type = forking

[Install]
WantedBy = multi-user.target
